# wazuh-odfe/config/etc/services.d/ossec-logs/run
#!/usr/bin/with-contenv sh

set -e

readonly _logfile='/var/ossec/logs/ossec.log'

# dumping ossec.log to standard output
exec tail -f "${_logfile}"

# vim:ft=sh
